# services/wiggle-sensors.timer
[Unit]
Description=Run wiggle-sensors.service every hour

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target