# services/wiggle-camera.service
[Unit]
Description=WiggleCamera service
After=network.target

[Service]
WorkingDirectory=/home/wiggle/WiggleBin/packages/venv/lib/python3.11/site-packages/wiggle_camera
ExecStart=/home/wiggle/WiggleBin/packages/venv/bin/python3 /home/wiggle/WiggleBin/packages/venv/lib/python3.11/site-packages/wiggle_camera/main.py --recording
Environment="PATH=/home/wiggle/WiggleBin/packages/venv/bin:/home/wiggle/WiggleBin/packages/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/usr/local/games:/usr/games"
Restart=always

[Install]
WantedBy=default.target